# supervisord.conf
[program:node]
command = node run_server.js
directory = /home/dotcloud/current